# data/sanoid/syncoid.timer
[Unit]
Description=Run Syncoid Hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
